# supervisor/incron.conf
[program:incron]
command=incrond -n
autostart=true
autorestart=true
